# source/supervisor-conf/kt_status.conf
[program:kt_status]
command = python3 /home/robot/robot/source/kt_status.py
directory = /home/robot/robot/source/
autostart = yes
autorestart = yes
user = root
stdout_logfile_maxbytes = 200MB
stderr_logfile_maxbytes = 200MB
